# aq_boss_linux/DEBIAN/postrm
#!/bin/bash
set -e
SERVICE_NAME="aq_boss_linux.service"
TARGET_DIR="/opt/boss-local"

case "$1" in
    remove)
        echo "Stopping and disabling ${SERVICE_NAME}..."
        systemctl stop ${SERVICE_NAME} 2>/dev/null || true
        systemctl disable ${SERVICE_NAME} 2>/dev/null || true
        ;;
    purge)
        echo "Completely removing ${TARGET_DIR}..."
        rm -rf "${TARGET_DIR}" 2>/dev/null || true
        ;;
esac

exit 0
