# aq_boss_linux/DEBIAN/prerm
#!/bin/bash
set -e

SERVICE_NAME="aq_boss_linux.service"
USB_RULE="/etc/udev/rules.d/80-aq-boss-usb.rules"
SYSTEMD_UNIT="/etc/systemd/system/aq_boss_usb_update@.service"
ACTION="$1"

echo "[AQ Boss Linux] Pre-removal step running (action: $ACTION)..."

# Skip cleanup on upgrade/deconfigure (self-update case)
if [ "$ACTION" = "upgrade" ] || [ "$ACTION" = "deconfigure" ]; then
    echo "Upgrade/deconfigure detected — skipping service and rule removal."
    exit 0
fi

# Stop and disable main service
if pidof systemd >/dev/null 2>&1 && command -v systemctl >/dev/null 2>&1; then
    echo "Stopping and disabling ${SERVICE_NAME}..."
    systemctl stop "${SERVICE_NAME}" 2>/dev/null || true
    systemctl disable "${SERVICE_NAME}" 2>/dev/null || true

    # Stop any active USB updater instances
    echo "Stopping AQ Boss USB updater instances..."
    systemctl list-units --all | grep -E "aq_boss_usb_update@" | awk '{print $1}' | while read -r svc; do
        systemctl stop "$svc" 2>/dev/null || true
    done

    systemctl daemon-reload || true
fi

echo "[AQ Boss Linux] Pre-removal step complete."
exit 0
