# aq-boss-linux/DEBIAN/postinst
#!/bin/bash
set -e

SERVICE_NAME="aq-boss-linux.service"
USB_RULE="/etc/udev/rules.d/80-aq-boss-usb.rules"
ACTION="$1"
PREVIOUS_VERSION="$2"

has_systemd() {
    pidof systemd >/dev/null 2>&1 && command -v systemctl >/dev/null 2>&1
}

is_upgrade=false
if [ "$ACTION" = "configure" ] && [ -n "$PREVIOUS_VERSION" ]; then
    is_upgrade=true
fi

echo "[AQ Boss Linux] Running post-install setup (action: $ACTION)..."

# Step 1: Reload udev only on fresh installs when rules are present
if ! $is_upgrade; then
    if [ -f "$USB_RULE" ] && command -v udevadm >/dev/null 2>&1; then
        echo "Reloading udev rules..."
        udevadm control --reload-rules || true
        udevadm trigger || true
    fi
fi

# Step 2: Reload systemd daemon whenever available
if has_systemd; then
    echo "Reloading systemd daemon..."
    systemctl daemon-reload || true
fi

# Step 3: Set permissions for updater script (always safe)
if [ -f "/opt/aq-boss-linux/updater/usb_update_check.py" ]; then
    chmod 755 /opt/aq-boss-linux/updater/usb_update_check.py
    chown root:root /opt/aq-boss-linux/updater/usb_update_check.py
fi

if [ -f "/opt/aq-boss-linux/updater/aq-boss-linux.py" ]; then
    chmod 755 /opt/aq-boss-linux/updater/aq-boss-linux.py
    chown root:root /opt/aq-boss-linux/updater/aq-boss-linux.py
fi

# Step 4: Ensure the main service is enabled and running
if has_systemd; then
    if [ -x "/opt/aq-boss-linux/updater/automount.sh" ]; then
        echo "Ensuring automount is enabled for the active desktop session..."
        if /opt/aq-boss-linux/updater/automount.sh --enable; then
            echo "Automount enabled."
        else
            echo "Warning: failed to enable automount automatically — enable USB automount manually to allow USB updates."
        fi
    fi
    if $is_upgrade; then
        echo "Restarting $SERVICE_NAME after upgrade..."
        systemctl restart "$SERVICE_NAME" >/dev/null 2>&1 || true
    else
        echo "Enabling and starting $SERVICE_NAME..."
        systemctl enable --now "$SERVICE_NAME" >/dev/null 2>&1 || true
    fi
fi

echo "[AQ Boss Linux] Post-install setup complete."
exit 0
