# aq_boss_linux/etc/udev/rules.d/80-aq-boss-usb.rules
ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_TYPE}=="vfat|ntfs|ext4", RUN+="/opt/boss-local/updater/usb_update_check.py"
